# root/etc/periodic/monthly/certs
#!/bin/sh
acme-client "$SSL_CERT_DOMAIN";
