# root/etc/periodic/monthly/certs
#!/bin/sh
acme-client -vF -c /certs $SSL_CERT_DOMAIN
