# root/etc/periodic/monthly/certs
#!/bin/sh
acme-client -vF $SSL_CERT_DOMAIN
